# cognihire-frontend.service
[Unit]
Description=CogniHire Next.js Frontend
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/cognihire
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
Environment=PORT=3000
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target